# stick-buddy.service
[Unit]
Description=Hermes Stick Buddy — VPS aggregation server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/home/linuxuser/.hermes/hermes-agent/venv/bin/python3 /home/linuxuser/hermes-stick-buddy/server/app.py
WorkingDirectory=/home/linuxuser/hermes-stick-buddy/server
EnvironmentFile=/home/linuxuser/hermes-stick-buddy/server/.env
# Set STICK_BUDDY_TOKEN in the .env file, never commit it here
Restart=always
RestartSec=3
User=linuxuser

[Install]
WantedBy=default.target